# tools/build/podman/s6/s6-rc.d/lanraragi/run
#!/bin/sh
cd /lanraragi/
export HOME=/lanraragi
exec perl /lanraragi/script/launcher.pl -f /lanraragi/script/lanraragi
